# ansible/roles/bme280tozabbix/files/bme280tozabbix.timer
[Unit]
Description=get temperature/humidity/pressure by bme280 and send to zabbix

[Timer]
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
